# noma_ai.service
[Unit]
Description=NOMA AI Skin Analysis System
After=multi-user.target
# Removed graphical.target dependencies – we boot to console

[Service]
Type=simple
User=havil
Group=havil
WorkingDirectory=/home/havil/noma_ai
ExecStart=/home/havil/noma_ai/start_noma.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=noma_ai
TimeoutStopSec=30
KillMode=process

# Environment is now set inside start_noma.sh – no DISPLAY/XAUTHORITY needed

[Install]
WantedBy=multi-user.target
